# infra/systemd/jogak-worker.service
[Unit]
Description=JOGAK GPU 7 worker
After=network-online.target redis.service
Wants=network-online.target

[Service]
Type=simple
User=guest
WorkingDirectory=/SSD/guest/chojoonghui/JOGAK
EnvironmentFile=/SSD/guest/chojoonghui/JOGAK/.env
Environment=CUDA_VISIBLE_DEVICES=7
Environment=NVIDIA_VISIBLE_DEVICES=7
Environment=PYTORCH_CUDA_ALLOC_CONF=expandable_segments:True
ExecStart=/home/guest/anaconda3/envs/cjh_jogak/bin/python worker/run_worker.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
